# deploy/investcool.service
[Unit]
Description=InvestCool Gunicorn Service
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/investcool
ExecStart=/opt/investcool/venv/bin/gunicorn -c deploy/gunicorn.conf.py run:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
